# conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=service nginx start
autorestart=false
redirect_stderr=true